# supervisord.conf
[supervisord]
nodaemon=true

[program:web]
command=/usr/sbin/nginx -g "daemon off;"
startsecs=0

[program:node_exporter]
command=/usr/local/bin/node_exporter
startsecs=0
